# systemd/.config/systemd/user/baterry-checker.timer
[Unit]
Description=Battery checker timer

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
