# systemd/.config/systemd/user/wallpaper.service
[Unit]
Description=wallpaper
[Service]
ExecStart=%h/bin/wallpaper
Restart=on-failure
[Install]
WantedBy=default.target
