# systemd/.config/systemd/user/wallpaper.timer
[Unit]
Description=wallpaper timer

[Timer]
OnBootSec=15sec
OnCalendar=hourly

[Install]
WantedBy=timers.target
